# builds/swi/all/shared/etc/udev/rules.d/60-block.rules
SUBSYSTEM=="block", RUN+="/sbin/initblockdev $kernel $env{ACTION}"
