# builds/swi/all/shared/etc/udev/rules.d/60-net.rules
SUBSYSTEM=="net", RUN+="/sbin/initnetdev $kernel $env{ACTION}"
