# rootconf/default/etc/init.d/rc.local
#!/bin/sh

#
# If present run /etc/rc.local for platform specific initialization
#

PATH=/usr/bin:/usr/sbin:/bin:/sbin

[ -x /etc/rc.local ] && /etc/rc.local $1

