# components/any/ofdpa-2.0-closed/deb/debuild/debian/ofdpa-2.0-closed.postinst
/sbin/ldconfig
/usr/sbin/update-rc.d ofdpa defaults
/usr/sbin/service ofdpa start



