# components/any/ofdpa-2.0-closed/deb/debuild/debian/ofdpa-2.0-closed.postinst
# for compat with older pre-production switches
ln -s /usr/sbin/ofdpa-powerpc-accton-as5710-54x-r0 /usr/sbin/ofdpa-powerpc-as5710-54x-r0b
ln -s /usr/sbin/ofdpa-powerpc-accton-as5710-54x-r0 /usr/sbin/ofdpa-powerpc-as5710-54x-r0a
/sbin/ldconfig
/usr/sbin/update-rc.d ofdpa defaults
/usr/sbin/service ofdpa start


