# demo/os/x86_64/etc/rc.local
#!/bin/sh

. /lib/onie/functions

# Mount demo filesystem
demo_mnt="/boot"

mkdir -p $demo_mnt || {
    echo "Error: Unable to create demo file system mount point: $demo_mnt"
    exit 1
}

mount -t ext4 -o defaults LABEL=ONIE-DEMO $demo_mnt || {
    echo "Error: Unable to mount $demo_dev on $demo_mnt"
    exit 1
}

# Local Variables:
# mode: shell-script
# eval: (sh-set-shell "/bin/sh" t nil)
# End:
