# machine/dell/dellemc_s4100_c2338/rootconf/sysroot-etc/inittab
# /etc/inittab for onie

#  Copyright (C) 2013 Curt Brune <curt@cumulusnetworks.com>
#
#  SPDX-License-Identifier:     GPL-2.0

::sysinit:/etc/init.d/rc S
::wait:/etc/init.d/rc 3

# Start an "askfirst" shell on the console (whatever that may be)\n"
::askfirst:-/bin/onie-console
ttyS1:12345:respawn:/sbin/getty 115200 ttyS1 vt102

::shutdown:/etc/init.d/rc 6
::restart:/sbin/init

