# machine/as5712_54x/rootconf/sysroot-etc/inittab
# /etc/inittab for onie

::sysinit:/etc/init.d/rc S
::wait:/etc/init.d/rc 3

# Start an "askfirst" shell on the console (whatever that may be)\n"
ttyS1::askfirst:-/bin/onie-console

::shutdown:/bin/umount -a -r
::restart:/sbin/init

